# microservices/metrics-manager/supervisord.conf
# Copyright (C) 2025-2026 Intel Corporation
# SPDX-License-Identifier: Apache-2.0

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
# The socket is created in /tmp with 0700 perms and is only reachable from
# inside the container, but supervisord still emits a CRIT line on every
# start unless an explicit username/password is configured. Using a
# non-secret placeholder pair silences the noise without weakening the
# actual access controls (Unix permissions remain the real gate).
username=supervisor
password=local-only-not-a-secret

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username=supervisor
password=local-only-not-a-secret

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
user=root
loglevel=warn
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[include]
files=/etc/supervisor/conf.d/*.conf

# Metrics Manager - FastAPI application
# Depends on Telegraf at runtime: it scrapes Telegraf's Prometheus endpoint
# (:9273) for the SSE stream and pushes custom metrics to Telegraf's HTTP
# listener (:8186). Telegraf is started first (lower priority) so the API has
# data to serve as soon as it is reachable; the API tolerates Telegraf being
# briefly unavailable, so no explicit readiness wait is needed.
[program:metrics-manager]
command=/opt/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port %(ENV_METRICS_PORT)s
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
environment=PYTHONPATH="/app",METRICS_PORT="%(ENV_METRICS_PORT)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

# Telegraf - System metrics collection
# Started before metrics-manager so its Prometheus endpoint (:9273) and HTTP
# listener (:8186) are reachable as soon as the API begins serving requests.
[program:telegraf]
command=telegraf --config %(ENV_TELEGRAF_CONFIG_PATH)s
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# qmassa - GPU metrics collection (Intel GPUs)
[program:qmassa]
command=qmassa --ms-interval 1000 --no-tui --to-json /app/qmassa.fifo
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
